# hw5/5.cron
55 23 * * * ./daily1.sh
0 8,17 * * * ./daily2.sh
30 8-18/1 * * 1-5 ./hourly.sh #есть сомнения
0 */4 * * * ./quatre.sh
